# deploy/systemd/secretariasistema.service
[Unit]
Description=Gestão Secretaria Backend API
Documentation=https://github.com/suaorg/gestao-secretaria
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/var/www/gestao-secretaria/backend
ExecStart=/usr/bin/node src/server.js
Restart=on-failure
RestartSec=10
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=gestao-secretaria

# Variáveis de ambiente (ou use EnvironmentFile=)
EnvironmentFile=/var/www/gestao-secretaria/backend/.env

# Garante que os diretórios necessários existam ANTES do namespace ser montado
ExecStartPre=/bin/mkdir -p /var/www/gestao-secretaria/backend/data
ExecStartPre=/bin/mkdir -p /var/www/gestao-secretaria/backend/uploads
ExecStartPre=/bin/mkdir -p /var/www/gestao-secretaria/backend/uploads/fotos
ExecStartPre=/bin/mkdir -p /var/www/gestao-secretaria/backend/uploads/logos
ExecStartPre=/bin/mkdir -p /var/www/gestao-secretaria/backend/uploads/carteiras
ExecStartPre=/bin/chown -R www-data:www-data /var/www/gestao-secretaria/backend/data
ExecStartPre=/bin/chown -R www-data:www-data /var/www/gestao-secretaria/backend/uploads

# Limites de segurança
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/www/gestao-secretaria/backend/data
ReadWritePaths=/var/www/gestao-secretaria/backend/uploads
CapabilityBoundingSet=
LimitNOFILE=65535
MemoryHigh=512M
MemoryMax=768M

[Install]
WantedBy=multi-user.target
