# deploy/systemd/secretariasistema.service
[Unit]
Description=SecretariaSistema Backend API
Documentation=https://github.com/suaorg/secretariasistema
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/var/www/secretariasistema/backend
ExecStart=/usr/bin/node src/server.js
Restart=on-failure
RestartSec=10
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=secretariasistema

# Variáveis de ambiente (ou use EnvironmentFile=)
EnvironmentFile=/var/www/secretariasistema/backend/.env

# Limites de segurança
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/www/secretariasistema/backend/data
ReadWritePaths=/var/www/secretariasistema/backend/uploads
CapabilityBoundingSet=
LimitNOFILE=65535
MemoryHigh=512M
MemoryMax=768M

[Install]
WantedBy=multi-user.target
